# device/etc/systemd/system/swupdate.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/local/bin/swupdate \
  -f /etc/swupdate-main.cfg \
  -k /etc/swupdate/pubkey.pem \
  -w "-p 8080"
